# supervisor/conf.d/apache.conf
[program:apache2]
command=%(ENV_APACHE_COMMAND)s
autostart=true
autorestart=true
stdout_logfile=/var/log/apache2/supervisor.log
stderr_logfile=/var/log/apache2/supervisor_err.log
